# alpine/etc/local.d/apps.start
#freq=`fw_printenv -n hw_rev | awk -F _ '{print $2}' | awk -F - '{print $1}'`
cat /media/mmcblk0p1/sdr_receiver_kiwi.bit > /dev/xdevcfg
cp /media/mmcblk0p1/kiwi.bin /root/kiwi.bin
cp -Rf /media/mmcblk0p1/config /root/config
cd /root
/root/kiwi.bin &
